# update-motd.d/05-greeting
#!/bin/sh
GREEN='\033[0;32m'
BRIGHT_GREEN='\033[1;32m'
NC='\033[0m' # No Color

echo "${GREEN}"
echo "‚ĖĎ‚ĖĎ‚ĖĎ‚ĖĎ‚ĖĎ‚ĖĎ‚ĖĎ‚ĖĎ‚ĖĎ‚ĖĎ‚ĖĎ‚ĖĎ‚ĖĎ‚ĖĎ‚ĖĎ‚ĖĎ‚ĖĎ‚ĖĎ‚ĖĎ‚ĖĎ‚ĖĎ‚ĖĎ‚ĖĎ‚ĖĎ‚ĖĎ‚ĖĎ‚ĖĎ‚ĖĎ‚ĖĎ‚ĖĎ‚ĖĎ‚ĖĎ‚ĖĎ‚ĖĎ‚ĖĎ‚ĖĎ‚ĖĎ‚ĖĎ‚ĖĎ‚ĖĎ‚ĖĎ‚ĖĎ‚ĖĎ‚ĖĎ‚ĖĎ‚ĖĎ‚ĖĎ‚ĖĎ‚ĖĎ‚ĖĎ‚ĖĎ‚ĖĎ‚ĖĎ‚ĖĎ‚ĖĎ‚ĖĎ‚ĖĎ‚ĖĎ‚ĖĎ‚ĖĎ‚ĖĎ‚ĖĎ"
echo "‚ĖĎ                                                            ‚ĖĎ"
# Note: Using cat << "EOF" to preserve those Unicode diagonal lines
cat << "EOF"
‚ĖĎ       ___        __      _  _         _    _               ‚ĖĎ
‚ĖĎ      ‚ēĪ _ ‚ē≤__  __‚ēĪ _‚Ēā_ __‚Ēā ‚Ēā‚Ēā ‚Ēā  _ __ ‚Ēā ‚Ēā _(_) ___          ‚ĖĎ
‚ĖĎ     ‚Ēā ‚Ēā ‚Ēā ‚ē≤ ‚ē≤‚ēĪ ‚ēĪ ‚Ēā_‚Ēā '__‚Ēā ‚Ēā‚Ēā ‚Ēā_‚Ēā '_ ‚ē≤‚Ēā ‚Ēā‚ēĪ ‚ēĪ ‚Ēā‚ēĪ _ ‚ē≤         ‚ĖĎ
‚ĖĎ     ‚Ēā ‚Ēā_‚Ēā ‚Ēā>  <‚Ēā  _‚Ēā ‚Ēā  ‚Ēā__   _‚Ēā ‚Ēā ‚Ēā ‚Ēā   <‚Ēā ‚Ēā  __‚ēĪ         ‚ĖĎ
‚ĖĎ      ‚ē≤___‚ēĪ‚ēĪ_‚ēĪ‚ē≤_‚ē≤_‚Ēā ‚Ēā_‚Ēā     ‚Ēā_‚Ēā ‚Ēā_‚Ēā ‚Ēā_‚Ēā_‚Ēā‚ē≤_‚ē≤_‚Ēā‚ē≤___‚Ēā         ‚ĖĎ
EOF
echo "‚ĖĎ                                                            ‚ĖĎ"
echo "‚ĖĎ                      [ ACCESS GRANTED ]                    ‚ĖĎ"
echo "‚ĖĎ                  Welcome back, ${BRIGHT_GREEN}0xfr4nkie!${GREEN}                  ‚ĖĎ"
echo "‚ĖĎ                                                            ‚ĖĎ"
echo "‚ĖĎ‚ĖĎ‚ĖĎ‚ĖĎ‚ĖĎ‚ĖĎ‚ĖĎ‚ĖĎ‚ĖĎ‚ĖĎ‚ĖĎ‚ĖĎ‚ĖĎ‚ĖĎ‚ĖĎ‚ĖĎ‚ĖĎ‚ĖĎ‚ĖĎ‚ĖĎ‚ĖĎ‚ĖĎ‚ĖĎ‚ĖĎ‚ĖĎ‚ĖĎ‚ĖĎ‚ĖĎ‚ĖĎ‚ĖĎ‚ĖĎ‚ĖĎ‚ĖĎ‚ĖĎ‚ĖĎ‚ĖĎ‚ĖĎ‚ĖĎ‚ĖĎ‚ĖĎ‚ĖĎ‚ĖĎ‚ĖĎ‚ĖĎ‚ĖĎ‚ĖĎ‚ĖĎ‚ĖĎ‚ĖĎ‚ĖĎ‚ĖĎ‚ĖĎ‚ĖĎ‚ĖĎ‚ĖĎ‚ĖĎ‚ĖĎ‚ĖĎ‚ĖĎ‚ĖĎ‚ĖĎ‚ĖĎ"
echo "${NC}"